# debian/perfsonar-elmond.postinst
#!/bin/sh
# postinst script for perfsonar-archive

# enable and start elmond
systemctl enable perfsonar-elmond.service
systemctl start perfsonar-elmond.service

#DEBHELPER#
exit 0
